# ansible_collections/oszi/general/roles/podman/templates/etc/systemd/system/podman-image-gc.timer.j2
[Unit]
Description=Podman Image Garbage Collector Timer

[Timer]
OnCalendar={{ podman_image_gc_on_calendar | default('daily') }}

[Install]
WantedBy=timers.target
