# ansible_collections/oszi/general/roles/podman/templates/etc/systemd/system/podman-image-gc.service.j2
[Unit]
Description=Podman Image Garbage Collector

[Service]
EnvironmentFile=-/etc/environment
Environment=PODMAN_RUNTIME={{ podman_runtime }}
ExecStart=/bin/sh -c "${PODMAN_RUNTIME} images -qf dangling=true | /usr/bin/xargs -r ${PODMAN_RUNTIME} rmi"
